# files/profile.d/50-enhancements.sh
alias ls='ls --color -l'
LS_COLORS='di=1:fi=0:ln=96:pi=5:so=5:bd=5:cd=5:or=31:mi=0:ex=35:*.rpm=90'
export LS_COLORS